# nxsiran-bridge.service
[Unit]
Description=NxSiran Bridge Service (独立服务)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/NxSiran-Telegram-Bot
EnvironmentFile=/root/NxSiran-Telegram-Bot/.env
ExecStart=/usr/bin/python3 /root/NxSiran-Telegram-Bot/bridge_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=BRIDGE_PORT=8081

[Install]
WantedBy=multi-user.target
