# system/nxsiran-webhook.service
[Unit]
Description=LoveSupremacy Webhook Listener (GitHub Auto-Deploy + Bridge)
Documentation=https://github.com/Jasper6439/LoveSupremacy_Universe
After=network.target nx_siran.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/LoveSupremacy-Telegram-Bot
EnvironmentFile=/opt/LoveSupremacy-Telegram-Bot/.env
ExecStart=/usr/bin/python3 /opt/LoveSupremacy-Telegram-Bot/system/webhook_server.py
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lovesupremacy-webhook
TimeoutStopSec=30
KillSignal=SIGTERM

# 安全加固
Environment=PYTHONUNBUFFERED=1
Environment=WEBHOOK_PORT=8082
Environment=PROJECT_DIR=/opt/LoveSupremacy-Telegram-Bot
Environment=BOT_SERVICE=nx_siran.service
Environment=BRIDGE_TOKEN=nxsiran_bridge_2024

# 内存限制（e2-micro 1GB RAM）
MemoryMax=128M

[Install]
WantedBy=multi-user.target
