# system/nxsiran-watcher.service
[Unit]
Description=NxSiran Bot Auto Restart Watcher
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/NxSiran-Telegram-Bot
ExecStart=/opt/NxSiran-Telegram-Bot/auto_restart_watcher.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
