# system/nxsiran-webhook.service
[Unit]
Description=NxSiran Webhook Service (独立服务)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/NxSiran/NxSiran-Telegram-Bot
EnvironmentFile=/opt/NxSiran/NxSiran-Telegram-Bot/.env
ExecStart=/usr/bin/python3 /opt/NxSiran/NxSiran-Telegram-Bot/system/webhook_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=WEBHOOK_PORT=8082
Environment=PROJECT_DIR=/opt/NxSiran/NxSiran-Telegram-Bot
Environment=BOT_SERVICE=nxsiran-bot.service

[Install]
WantedBy=multi-user.target
